# worker/deploy/supervisord.conf
; Production worker process map. Brings the validated headless stack up in a
; RACE-FREE order so the uinput devices are bound by libinput WITHOUT an Xorg
; restart: the input daemons create the devices, a one-shot binds/retags them,
; and Xorg is gated on that completing (run-xorg.sh waits for /tmp/uinput-ready)
; — then Xfce, optional VNC for the one-time Steam login, and the worker agent
; (which performs the spectate bring-up: GUI Steam login, Dota, GUI spectate,
; FFmpeg). Run as PID 1 with:  dumb-init -- supervisord -c /etc/supervisord.worker.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; udev provider — before everything so libudev input enumeration has a runtime.
[program:udev]
command=/usr/local/bin/start-udev-auto.sh
user=root
autorestart=true
priority=1
stopsignal=INT
stdout_logfile=/tmp/udev.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; System dbus bus — steamwebhelper's CEF needs it.
[program:dbus]
command=/usr/local/bin/start-dbus.sh
user=root
autorestart=true
priority=5
stdout_logfile=/tmp/dbus.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; Input daemons own the kernel uinput devices (created BEFORE Xorg). They must
; exist when Xorg does its startup input enumeration, so they precede it.
[program:uinput_daemon]
command=/usr/local/bin/uinput_daemon.py /tmp/dota_uinput.fifo
user=worker
autorestart=true
priority=6
stdout_logfile=/tmp/uinput_daemon.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:vnc_input_daemon]
command=/usr/local/bin/vnc_input_daemon.py /tmp/vnc_input.fifo
user=worker
environment=VNC_SCREEN="%(ENV_SCREEN)s"
autorestart=true
priority=6
stdout_logfile=/tmp/vnc_input_daemon.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; One-shot: wait for the 3 devices, ensure their /dev/input nodes, retag the
; udev db (dumb-udev mis-tags them JOYSTICK; libinput then refuses them), then
; touch /tmp/uinput-ready so Xorg may start. Runs as root (mknod + /run/udev).
[program:input_bind]
command=/usr/local/bin/input-bind.sh
user=root
autostart=true
autorestart=false
startsecs=0
priority=7
stdout_logfile=/tmp/input-bind.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; Real NVIDIA Xorg on :99, gated on input-bind (run-xorg.sh waits for the ready flag).
[program:xorg]
command=/usr/local/bin/run-xorg.sh
user=root
autorestart=true
priority=10
stdout_logfile=/tmp/xorg.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; Full Xfce4 session as the worker user (its compositor is what makes Steam's CEF render).
[program:xfce]
command=/usr/local/bin/start-xfce.sh
user=worker
environment=HOME="/fard/steam/steamhome",USER="worker",DISPLAY=":99"
autorestart=true
priority=20
startsecs=5
stdout_logfile=/tmp/xfce.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; x11vnc + noVNC — gated by ENABLE_VNC; set true once over an SSH tunnel for the
; one-time interactive Steam login, false thereafter (silent auto-login).
[program:x11vnc]
command=/usr/local/bin/start-vnc.sh
user=worker
environment=HOME="/fard/steam/steamhome",DISPLAY=":99",ENABLE_VNC="%(ENV_ENABLE_VNC)s"
autorestart=true
priority=30
stdout_logfile=/tmp/x11vnc.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:novnc]
command=/usr/local/bin/start-novnc.sh
user=worker
environment=HOME="/fard/steam/steamhome",ENABLE_VNC="%(ENV_ENABLE_VNC)s"
autorestart=true
priority=40
stdout_logfile=/tmp/novnc.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; The worker agent — waits for Xorg + the bound input devices, then connects to
; the control plane and drives the spectate lifecycle. Env from the container.
[program:worker]
command=/usr/local/bin/run-worker.sh
user=worker
environment=HOME="/fard/steam/steamhome",DISPLAY=":99",WORKER_DOTA_BRINGUP="%(ENV_WORKER_DOTA_BRINGUP)s",CONTROL_PLANE_ADDR="%(ENV_CONTROL_PLANE_ADDR)s",MEDIAMTX_SRT_HOST="%(ENV_MEDIAMTX_SRT_HOST)s",MEDIAMTX_SRT_PORT="%(ENV_MEDIAMTX_SRT_PORT)s",SCREEN="%(ENV_SCREEN)s"
autorestart=true
priority=50
startsecs=10
stdout_logfile=/tmp/worker.log
stdout_logfile_maxbytes=0
redirect_stderr=true
