# scripts/validation/supervisord.conf
; Supervisor config for the V4/V5 headless Steam desktop (replicates docker-steam-headless's
; process management). Brings up, in order: the real NVIDIA Xorg :99 (root, GPU) -> a full
; Xfce4 session (worker uid; xfwm4 + compositor, which is what makes Steam's CEF login popup
; render) -> x11vnc + noVNC for the operator's one-time interactive login.
;
; Steam is NOT started here — V4 launches it via `docker exec` so it can drive the login phase
; (interactive, ENABLE_VNC=true) and the auto-login proof phase (ENABLE_VNC=false) separately.
;
; Run as PID 1 with:  dumb-init -- supervisord -c /etc/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

; Local socket so `supervisorctl status` works for inspecting program state.
[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; System dbus bus — steamwebhelper's CEF needs it (otherwise it spams connection errors).
[program:dbus]
command=/usr/local/bin/start-dbus.sh
user=root
autorestart=true
priority=5
stdout_logfile=/tmp/dbus.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; Xorg runs as root (needs the seat/GPU); -ac disables X access control so the non-root
; worker user can open the display.
[program:xorg]
command=Xorg :99 -config /etc/X11/xorg.conf -noreset -ac
user=root
autorestart=true
priority=10
stdout_logfile=/tmp/xorg.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; Full Xfce4 desktop session as the worker user, with its own session bus.
[program:xfce]
command=/usr/local/bin/start-xfce.sh
user=worker
environment=HOME="/fard/steam/steamhome",USER="worker",DISPLAY=":99"
autorestart=true
priority=20
startsecs=5
stdout_logfile=/tmp/xfce.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; x11vnc — gated by ENABLE_VNC (the wrapper idles when false, for the headless auto-login proof).
[program:x11vnc]
command=/usr/local/bin/start-vnc.sh
user=worker
environment=HOME="/fard/steam/steamhome",DISPLAY=":99",ENABLE_VNC="%(ENV_ENABLE_VNC)s"
autorestart=true
priority=30
stdout_logfile=/tmp/x11vnc.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; noVNC (websockify) so the operator can reach the desktop from a browser over an SSH tunnel.
[program:novnc]
command=/usr/local/bin/start-novnc.sh
user=worker
environment=HOME="/fard/steam/steamhome",ENABLE_VNC="%(ENV_ENABLE_VNC)s"
autorestart=true
priority=40
stdout_logfile=/tmp/novnc.log
stdout_logfile_maxbytes=0
redirect_stderr=true
